# etc/skel/.config/systemd/user/polkit-gnome-authentication-agent.service
[Unit]
Description=Authentication agent
Documentation=https://wiki.archlinux.org/title/Polkit

PartOf=hyprland-session.target
After=hyprland-session.target

[Service]
Type=exec
ExecStart=/usr/lib/polkit-gnome/polkit-gnome-authentication-agent-1
ExecStop=pkill -f polkit-gnome-authentication-agent-1

[Install]
WantedBy=hyprland-session.target
